# timelapse.service
[Unit]
Description=Timelapse Camera Service
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /opt/timelapse/timelapse.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target